# conf.d/firefox.conf
[program:firefox]
command=firefox
autorestart=true
